# crontab
#crontab
0 3 * * * /usr/local/bin/python3 /backup.py clean >> /var/log/backup.log 2>&1
0 14 * * * /usr/local/bin/python3 /backup.py >> /var/log/backup.log 2>&1
#end of line
